# etc/services.d/syslog-ng/run
#!/usr/bin/execlineb
/usr/sbin/syslog-ng -F -f /etc/syslog-ng/syslog-ng.conf
